# cron.daily/pandora_db
#!/bin/bash
perl /usr/share/pandora_server/util/pandora_db.pl /etc/pandora/pandora_server.conf
